# systemd/osia-mcp-time-bridge.service
[Unit]
Description=OSIA Time MCP Bridge
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/osia-framework
ExecStart=/home/ubuntu/osia-framework/.venv/bin/python src/mcp_bridge.py --name time --command uvx --port 8004 --args mcp-server-time
Restart=always
RestartSec=5
EnvironmentFile=/home/ubuntu/osia-framework/.env
Environment=PYTHONPATH=/home/ubuntu/osia-framework

[Install]
WantedBy=multi-user.target
