# systemd/osia-research-worker.timer
[Unit]
Description=OSIA Research Worker — runs every 15 minutes
Requires=osia-research-worker.service

[Timer]
# Fire every 15 minutes to keep pace with queue growth
OnCalendar=*:0/15:00
AccuracySec=30s
Persistent=false

[Install]
WantedBy=timers.target
