# systemd/osia-daily-sitrep.service
[Unit]
Description=Trigger OSIA Daily SITREP
After=network.target osia-orchestrator.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/osia-framework
ExecStart=/home/ubuntu/osia-framework/.venv/bin/python -m src.cron.daily_sitrep
EnvironmentFile=/home/ubuntu/osia-framework/.env

[Install]
WantedBy=multi-user.target
