# systemd/osia-mcp-tavily-bridge.service
[Unit]
Description=OSIA Tavily MCP Bridge
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/osia-framework
ExecStart=/home/ubuntu/osia-framework/.venv/bin/python src/mcp_bridge.py --name tavily --command npx --port 8003 --args -y @tavily/mcp-server
Restart=always
RestartSec=5
EnvironmentFile=/home/ubuntu/osia-framework/.env
Environment=PYTHONPATH=/home/ubuntu/osia-framework

[Install]
WantedBy=multi-user.target
